# disease-api.service
[Unit]
Description=Disease Prediction API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/DocAI
Environment="PATH=/path/to/DocAI/myenv/bin"
ExecStart=/path/to/DocAI/myenv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
